# neomutt/.config/systemd/user/msmtpq.timer
[Unit]
Description=Msmtpq timer

[Timer]
OnBootSec=1m
OnUnitActiveSec=30s

[Install]
WantedBy=timers.target
