# systemd/.config/user/vdirsyncer.service
[Unit]
Description=Run vdirsyncer

[Service]
Type=simple
ExecStart=/usr/bin/vdirsyncer sync

[Install]
WantedBy=default.target
