# systemd/.config/systemd/user/cloud-audacis.timer
[Unit]
Description=Sync audacis cloud every minute

[Timer]
OnBootSec=1m
OnUnitActiveSec=2min

[Install]
WantedBy=timers.target
