# systemd/.config/systemd/user/linkpaper.service
[Unit]
Description=Link scientific paper to common location

[Service]
ExecStart=/home/phil/bin/linkpaper
StandardOutput=syslog
