# systemd/.config/systemd/user/lions.service
[Unit]
Description=Check Lions calendar

[Service]
ExecStart=/home/phil/bin/lions -m
StandardOutput=syslog
