# systemd/.config/systemd/user/mailnotify.service
[Unit]
Description=Wendy for notimail

[Service]
ExecStart=/usr/bin/wendy -m 128 -f "${MAIL}/new" notimail

[Install]
WantedBy=default.target
