# systemd/.config/systemd/user/mbsync.service
[Unit]
Description=Sync mails with mbsync
Documentation=man:mbsync(1)
TimeoutStartSec=30s

[Service]
Type=oneshot
ExecStart=/usr/bin/mbsync -Va
StandardOutput=syslog
StandardError=syslog
