# systemd/.config/systemd/user/mbsync.timer
[Unit]
Description=Run mbsync every minute

[Timer]
OnCalendar=minutely
Persistent=true
Unit=mbsync.service

[Install]
WantedBy=timers.target
