# systemd/.config/systemd/user/mpd.service
[Unit]
Description=Start mpd
Documentation=man:mpd(1)

[Service]
ExecStart=/usr/bin/mpd --no-daemon

[Install]
WantedBy=default.target
