# systemd/.config/systemd/user/mpdbar.service
[Unit]
Description=Start mpd bar
Requires=mpd.service

[Service]
ExecStart=/home/phil/bin/mpdbar

[Install]
WantedBy=default.target
