# systemd/.config/systemd/user/nextcloud.timer
[Unit]
Description=Sync nextcloud every minute

[Timer]
OnCalendar=*:0/15
Persistent=true
Unit=nextcloud.service

[Install]
WantedBy=timers.target
