# systemd/.config/systemd/user/redshift.service
[Unit]
Description=Start redshift
Documentation=man:redshift(1)

[Service]
ExecStart=/usr/bin/redshift -v

[Install]
WantedBy=default.target
