# systemd/.config/systemd/user/sandmann.timer
[Unit]
Description=Run sandmann daily

[Timer]
OnCalendar=*-*-* 18:30:00
Persistent=true
Unit=sandmann.service

[Install]
WantedBy=timers.target
