# systemd/.config/systemd/user/vdirsyncer.timer
[Unit]
Description=Run vdirsyncer daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
