# systemd/.config/systemd/user/vdirsyncer.service
[Unit]
Description=Run vdirsyncer

[Service]
ExecStart=/usr/bin/vdirsyncer sync

[Install]
WantedBy=default.target
