# systemd/.config/systemd/user/vdirsyncer.timer
[Unit]
Description=Sync vdirsyncer hourly

[Timer]
OnCalendar=hourly
Persistent=true
Unit=vdirsyncer.service

[Install]
WantedBy=timers.target
