# systemd/.config/systemd/user/vdirsyncer.timer
[Unit]
Description=Sync vdirsyncer every 5min

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
