# systemd/.config/systemd/user/vdirsyncer.service
[Unit]
Description=Sync vdirsyncer
Documentation=man:vdirsyncer(1)

[Service]
Type=oneshot
ExecStart=/usr/bin/vdirsyncer sync
